# srcs/mariadb/srcs/supervizord.conf
[supervisord]
nodaemon=true
user=root

[program:mysql]
command = sh "./init_db.sh"
